# vllm/thermal-guard/systemd/dcgm-exporter.service
[Unit]
Description=DCGM Exporter - NVIDIA GPU metrics exporter
Documentation=https://github.com/NVIDIA/dcgm-exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Start DCGM exporter with default counters
ExecStart=/usr/bin/dcgm-exporter -f /etc/dcgm-exporter/default-counters.csv

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dcgm-exporter

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
